# tvbutton.service
[Unit]
Description=tvbutton
After=network.target

[Service]
WorkingDirectory=/var/git/onepiecetv/
ExecStart=/usr/bin/python3 /var/git/onepiecetv/buttons.py
Restart=always

[Install]
WantedBy=multi-user.target
